# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

crontab /app/config/crontab.txt

exec crond -f;